# gnome-macos-remap.service
[Unit]
Description=GNOME → macOS Remap
After=graphical-session.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=pkexec /usr/local/bin/xremap %h/.config/gnome-macos-remap/config.yml

[Install]
WantedBy=graphical-session.target
